# systemd/.config/systemd/user/limit-2ram.slice
[Unit]
Description=Custom slice for limit resources usage

[Slice]
MemoryMax=2048M
MemorySwapMax=512M
