# xremap/.config/systemd/user/xremap.service
[Unit]
Description=Xremap tool daemon

[Service]
ExecStart=/home/alcides/.local/cargo/bin/xremap /home/alcides/.config/xremap/config.yml

[Install]
WantedBy=default.target
